# alacritty/systemd.service
[Unit]
Description=Per-monitor alacritty theme daemon

[Service]
Type=simple
ExecStart=%h/.local/bin/alacritty-theme-daemon
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
